# srcs/files/wolg.service
[Unit]
Description=Enable Wake-On-Lan
 
[Service]
Type=simple
ExecStart=/home/wolg
Restart=no
 
[Install]
WantedBy=multi-user.target